# conf.d/grass-node.conf
[program:grass-node]
command=python3 /app/grass-node_main.py
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0